# .config/autostart/screensaver-settings.desktop
[Desktop Entry]
Name=Screensaver
Comment=Set screensaver timeouts
Exec=xset s 300 dpms 420 0 0
