# obsidian/.config/obsidian/launchd/com.seiji.claude-session-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.seiji.claude-session-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>$HOME/.config/obsidian/scripts/claude-session-sync.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/claude-session-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-session-sync.error.log</string>

    <key>RunAtLoad</key>
    <false/>

</dict>
</plist>
